# services/php-fpm/run
#!/bin/bash

exec /usr/sbin/php-fpm7.0 -F -y /config/nginx/nginx-fpm.conf


